# masterref/global/etc/xinetd.d/pop3
service pop3
{
        socket_type             = stream
        protocol                = tcp
        wait                    = no
        user                    = root
        server                  = /usr/local/sbin/popper
        server_args             = qpopper -f /etc/qpopper/config
        disable                 = no
}

